# utils/unix/DebianInstaller/etc/init.d/adempiere
#!/bin/bash
#

### BEGIN INIT INFO
# Provides:		adempiere adempiere-3.6
# Required-Start:	postgresql
# Required-Stop:	postgresql
# Default-Start:	2 3 4 5
# Default-Stop:		1
# Short-Description:	ADempiere 3.6 server
# Description:		Provides ADempiere ERP-CRM Server startup and shutdown script. Requires PostgreSQL server.
# FileTarget:	/etc/init.d/adempiere
# FileOwner:	root.root
# FilePerms:	0755
#
# chkconfig:	2345 97 06
### END INIT INFO

# processname: adempiere
# Red Hat or SuSE config: /etc/sysconfig/adempiere
# Debian or Ubuntu config: /etc/default/adempiere

# Source function library
if [ -f /lib/lsb/init-functions ]
then
	. /lib/lsb/init-functions
elif [ -f /etc/init.d/functions ]
then
	. /etc/init.d/functions
fi

# Set path if path not set (if called from /etc/rc)
case $PATH in
    "") PATH=/bin:/usr/bin:/sbin:/etc
        export PATH ;;
esac

# initialization
# adjust these variables to your environment
ADEMPIERE_USER=adempiere
ADEMPIERE_HOME=/opt/Adempiere
JAVA_HOME=/usr/lib/jvm/java-6-openjdk
PATH=$PATH:${JAVA_HOME}/bin
SU=su
STOPMESSAGE="INFO.*Server\].*Shutdown complete" # Message when using java 6
export ADEMPIERE_HOME
export JAVA_HOME

if [ $(id -u) != "0" ]
then
    echo "You must be root to run the configure script.  Login as root and then run the configure script."
    exit 1
fi

RETVAL=0
ADEMPIERESTATUS=
MAXSECONDS=120 # max wait 2 minutes
SLEEPSECONDS=2
MAXITERATIONS=`expr $MAXSECONDS / $SLEEPSECONDS`

CONFIG_NAME=adempiere
CONFIGURATION="/etc/sysconfig/$CONFIG_NAME"
if [ -f /etc/redhat-release ]
then
    . /etc/init.d/functions

    init_status()
    {
	return 0
    }
    exit_status()
    {
	exit $?
    }
    success_status()
    {
	success
	echo
    }
    failure_status()
    {
	failure $?
	echo
    }

elif [ -f /etc/SuSE-release ]
then
    . /etc/rc.status

    init_status()
    {
	rc_reset
    }
    success_status()
    {
	echo "OK"
	return 0
    }
    failure_status()
    {
	echo "Failed"
	return 1
    }
    exit_status()
    {
	exit $?
    }

else
    if [ -d /etc/default ]
    then
        CONFIGURATION="/etc/default/$CONFIG_NAME"
    fi

    init_status()
    {
        return 0
    }

    success_status()
    {
        echo "OK"
        return 0
    }

    failure_status()
    {
        echo "Failed"
        return 0
    }

    exit_status()
    {
        exit $?
    }
fi

# Source configuration

[ -f "$CONFIGURATION" ] && . "$CONFIGURATION"

init_status

#
# if_fail()
#
# Evaluates return codes.  If 0, prints "OK", if 1, prints "Failed"
# and exits.  If 2, status is "already done" and nothing is printed.
# The rest of the functions in here all honor this convention.
#
if_fail()
{
    RC="$1"
    REASON="$2"
    if [ "$RC" = "0" ]
    then
        return
    elif [ "$RC" = "2" ]
    then
        return
    fi
    failure_status "${REASON}"
    exit 1
}

#
# write_sysconfig()
#
# Writes the system configuration
#
write_sysconfig()
{
	cat >"$CONFIGURATION" <<EOF
    
#This is a configuration file for automatic starting of the Adempiere
#Server at system startup. It is generated by running
#'/etc/init.d/adempiere configure'.Please use that method to modify this 
#file

# RUN_AT_STARTUP:'true' means to load the Server at system boot.
RUN_AT_STARTUP=${RUN_AT_STARTUP:-false}

# Configuration : Check whether configure has been done or not
CONFIGURE_RUN=${CONFIGURE_RUN}

EOF

    if [ $? != 0 ]
	then
		return 1
	fi
	return 0
}

# configure_perform()
# Create AdempiereEnv.properties
# RUN_silentsetup.sh
# utils/RUN_ImportAdempiere.sh
configure_perform()
{
    cp ${ADEMPIERE_HOME}/AdempiereEnvTemplate.properties ${ADEMPIERE_HOME}/AdempiereEnv.properties
    sed -i "s:^ADEMPIERE_HOME=.*:ADEMPIERE_HOME=${ADEMPIERE_HOME}:" ${ADEMPIERE_HOME}/AdempiereEnv.properties
    sed -i "s:^JAVA_HOME=.*:JAVA_HOME=${JAVA_HOME}:" ${ADEMPIERE_HOME}/AdempiereEnv.properties
    sed -i "s:^ADEMPIERE_DB_PASSWORD=.*:ADEMPIERE_DB_PASSWORD=${ADEMPIERE_DB_PASSWORD}:" ${ADEMPIERE_HOME}/AdempiereEnv.properties
    sed -i "s:^ADEMPIERE_DB_SYSTEM=.*:ADEMPIERE_DB_SYSTEM=${ADEMPIERE_DB_SYSTEM}:" ${ADEMPIERE_HOME}/AdempiereEnv.properties
    sed -i "s:^ADEMPIERE_WEB_PORT=.*:ADEMPIERE_WEB_PORT=${ADEMPIERE_WEB_PORT}:" ${ADEMPIERE_HOME}/AdempiereEnv.properties
    sed -i "s:^ADEMPIERE_SSL_PORT=.*:ADEMPIERE_SSL_PORT=${ADEMPIERE_SSL_PORT}:" ${ADEMPIERE_HOME}/AdempiereEnv.properties
    sed -i "s:^ADEMPIERE_APPS_SERVER=.*:ADEMPIERE_APPS_SERVER=0.0.0.0:" ${ADEMPIERE_HOME}/AdempiereEnv.properties
    sed -i "s:^ADEMPIERE_JAVA_TYPE=.*:ADEMPIERE_JAVA_TYPE=OpenJDK:" ${ADEMPIERE_HOME}/AdempiereEnv.properties
    sed -i "s:^ADEMPIERE_KEYSTORE=.*:ADEMPIERE_KEYSTORE=${ADEMPIERE_HOME}/keystore/myKeystore:" ${ADEMPIERE_HOME}/AdempiereEnv.properties
    chown adempiere:adempiere ${ADEMPIERE_HOME}/AdempiereEnv.properties
    chmod 600 ${ADEMPIERE_HOME}/AdempiereEnv.properties

    sed -i "s/:8080/:${ADEMPIERE_WEB_PORT}/g" /usr/share/applications/adempiere-homepage.desktop
    sed -i "s/:8443/:${ADEMPIERE_SSL_PORT}/g" /usr/share/applications/adempiere-webclient.desktop

    if [ -f ${ADEMPIERE_HOME}/RUN_silentsetup.sh ]  
    then
     	echo -n "Deploying ADempiere ERP server to jboss..."
       	$SU ${ADEMPIERE_USER} -c "cd ${ADEMPIERE_HOME}; ./RUN_silentsetup.sh" # > /dev/null 2>&1
	echo "Done"
    fi
    grep ${ADEMPIERE_WEB_PORT} ${ADEMPIERE_HOME}/jboss/server/adempiere/deploy/jboss-web.deployer/server.xml > /dev/null 2>&1
    if [ $? -ne 0 ]
    then
        echo "Error deploying ADempiere, please verify and try again"
        exit 1
    fi

    if [ -z "${ADEMPIERE_DB_SYSTEM}" ]
    then
        # recreate adempiere user from here as the RUN_ImportAdempiere requires it
        $SU postgres -c "dropdb -U postgres adempiere" # > /dev/null 2>&1
        $SU postgres -c "dropuser -U postgres adempiere" # > /dev/null 2>&1
        $SU postgres -c "psql -U postgres -c \"CREATE ROLE adempiere SUPERUSER LOGIN PASSWORD '${ADEMPIERE_DB_PASSWORD}'\"" # > /dev/null 2>&1
    fi

    if [ -f ${ADEMPIERE_HOME}/utils/RUN_ImportAdempiere.sh ]  
    then
        echo -n "Importing seed database..."
       	$SU ${ADEMPIERE_USER} -c "cd ${ADEMPIERE_HOME}/utils; ( echo "" | ./RUN_ImportAdempiere.sh )" # > /dev/null 2>&1
	echo "Done"
    fi
    $SU ${ADEMPIERE_USER} -c "export PGPASSWORD=${ADEMPIERE_DB_PASSWORD}; psql -d adempiere -U adempiere -h localhost -p 5432 -c ''" > /dev/null 2>&1
    if [ $? -ne 0 ]
    then
        echo "Cannot connect to adempiere database, please verify and try again"
        exit 1
    fi
    $SU ${ADEMPIERE_USER} -c "export PGPASSWORD=${ADEMPIERE_DB_PASSWORD}; psql -d adempiere -U adempiere -h localhost -p 5432 -c 'select count(*) from ad_system' 2>&1 | grep '1$'" > /dev/null 2>&1
    if [ $? -ne 0 ]
    then
        echo "Database not imported correctly, please verify and try again"
        exit 1
    fi

    chmod 700 ${ADEMPIERE_HOME}/utils/myEnvironment.*

    echo "Installation Completed Successfully."

    return 0
}

checkportused()
{
port=`netstat -n --tcp --listen | grep :$1 | awk '{print $4}' | cut -d':' -f2`
if [ "$port" = "$1" ]
then
    return 0  # Used
else
    return 1  # Not Used
fi
}

#
#configure_ask()
#
# Ask configuration questions,setting the variables.
#
configure_ask()
{
	cat <<EOF

ADempiere ERP Server Configuration
-------------------------------------------------
This will configure on-boot properties of ADempiere ERP Server.
The following questions will determine whether the database should 
be starting upon system boot, the ports it will use, and the passwords that 
will be used for database accounts.  Press <Enter> to accept the defaults. 
Ctrl-C will abort.

EOF
    # TODO: ask for certificate data (or integrate with openssl certificate)
    # TODO: ask for mail server settings
    # TODO: ask for postgres port (and perhaps machine if it's remote)

    # check JAVA_HOME
    if [ ! -x ${JAVA_HOME}/bin/jarsigner ]
    then
        echo ${JAVA_HOME}/bin/jarsigner is required.  Verify installation of sun-java6-jdk package.
        exit 1
    fi

    # Check not configurable JBoss ports (i.e. 1098, 1099, 3873, 4444, 4445, 4446, 8009, 8083, 8090, 8092, 8093)
    for JBOSSPORT in 1098 1099 3873 4444 4445 4446 8009 8083 8090 8092 8093
    do
        if checkportused "$JBOSSPORT"
        then
            echo Port $JBOSSPORT is required by ADempiere ERP server and it appears to be in use by another application.\
            Please fix the problem and try again.
            exit 1
        fi
    done

    #get the http port value
    while :
    do
	while [ 1 ]
	do
            echo -n Specify the HTTP port that will be used for ADempiere server [8080]:
            read LINE
            if [ -z $LINE ]
            then
                LINE=8080
            fi
            if checkportused "$LINE"
            then
                echo Port $port appears to be in use by another application.\
                Please specify a different port.
            else
                break;
            fi
        done

	case "$LINE" in
	"")
            break
            ;;
        *[^0-9]*)
            echo "Invalid http port: $LINE"
            ;;
        *)
            ADEMPIERE_WEB_PORT=$LINE
            break
            ;;
	esac
    done
    
    #get the https port value
    while :
    do
	while [ 1 ]
	do
            echo -n Specify the HTTPS port that will be used for ADempiere server [8443]:
            read LINE
            if [ -z $LINE ]
            then
                LINE=8443
            fi
            if checkportused "$LINE"
            then
                echo Port $port appears to be in use by another application.\
                Please specify a different port.
            else
                break;
            fi
        done

	case "$LINE" in
	"")
            break
            ;;
        *[^0-9]*)
            echo "Invalid https port: $LINE"
            ;;
        *)
            ADEMPIERE_SSL_PORT=$LINE
            break
            ;;
	esac
    done

    #get the adempiere database password
    while :
    do
        echo -n "Specify a password to be used for adempiere database account:"
	while [ 1 ]
	do
	    stty -echo > /dev/null 2>&1
	    temp=`echo $IFS`
            export IFS="\n"	
	    while [ 1 ]
  	    do
	        read LINE
                while [ -z "$LINE" ]
	        do
		    echo
		    echo -n "Password can't be null. Enter password:"  
	            read LINE
	        done

	        result=`expr index "$LINE" [\'\"]`
	        if [ $result != 0 ];
	        then
             	    echo 
	            echo -n "The password you entered contains invalid characters. Enter password:"	
	        else
		    break
	        fi
	    done
	    echo
	    echo -n "Confirm the password:"
	    read LINE1
	    echo
            if [ "$LINE" != "$LINE1" ];
	    then
		echo    
		echo -n "Passwords do not match.  Enter the password:"
	    else
		break
	    fi
	done
        stty echo > /dev/null 2>&1
        ADEMPIERE_DB_PASSWORD=$LINE
	export IFS=$temp
        break;
    done

    #get the postgres database password
    while :
    do
        echo -n "Specify the password of the user postgres on postgres database
(if empty then local connection will be tried):"
	while [ 1 ]
	do
	    stty -echo > /dev/null 2>&1
	    temp=`echo $IFS`
            export IFS="\n"	
	    while [ 1 ]
  	    do
	        read LINE
                if [ -z "$LINE" ]
	        then
		    break
	        fi

	        result=`expr index "$LINE" [\'\"]`
	        if [ $result != 0 ];
	        then
             	    echo 
	            echo -n "The password you entered contains invalid characters. Enter password:"
	        else
		    break
	        fi
	    done
	    echo
            if [ -z "$LINE" ]
            then
                # Empty postgres password verify connection with local postgres user
                $SU postgres -c "psql -U postgres -c ''"
                if [ $? -eq 0 ]
                then
                    break
                fi
                echo -n "Could not connect locally with user postgres.
NOTE: You can verify pg_hba.conf to check if this line is enabled
local   all         postgres                          ident
Enter postgres password:"
            else
                PGPASSWORD=$LINE
                export PGPASSWORD
                psql -U postgres -d template1 -h localhost -p 5432 -c ''
                if [ $? -eq 0 ]
                then
                    break
                fi
                echo -n "Could not connect with user postgres to database template with such password.
NOTE: You can check the file pg_hba.conf to validate that connection is allowed,
      or check that postgresql server allows tcp/ip connections on port 5432
      or check that user postgres database password is correctly set.
Enter postgres password:"
            fi
	done
        stty echo > /dev/null 2>&1
        ADEMPIERE_DB_SYSTEM=$LINE
	export IFS=$temp
        break;
    done

    while :
    do
        if [ "${RUN_AT_STARTUP}" = "true" ]
        then
            CUR=y
        else
            CUR=n
        fi
	echo
        echo -n "Do you want ADempiere ERP Server to be started on boot (y/n) [y]:"
        read LINE
	if [ -z $LINE ]
	then
	    RUN_AT_STARTUP=true
	fi
        echo
        case "$LINE" in
        "") 
            break
            ;;
        y|Y)
            RUN_AT_STARTUP=true
            break
            ;;
        n|N)
            RUN_AT_STARTUP=false
            break
            ;;
        *)
            echo "Invalid response: $LINE " >&2
            break
        esac
    done
}

getadempierestatus() {
    ADEMPIERESTATUSSTRING=$(ps ax | grep -v grep | grep ${ADEMPIERE_HOME})
    echo $ADEMPIERESTATUSSTRING | grep -q ${ADEMPIERE_HOME}
    ADEMPIERESTATUS=$?
}

configure()
{
	if [ "${CONFIGURE_RUN}" = "true" ]
	then
		echo "Adempiere is already configured"
	exit 1
	fi
	configure_ask
	configure_perform
	CONFIGURE_RUN=true
	write_sysconfig
	echo To access the ADempiere Server Home Page, start the server with '/etc/init.d/adempiere start'
        echo and then go to \"http://127.0.0.1:${ADEMPIERE_WEB_PORT}/admin\"
        
}

start () {

    if [ "${CONFIGURE_RUN}" != "true" ]		
    then
	echo "ADempiere is not configured.  You must run
'/etc/init.d/adempiere configure' as the root user to configure the server."
	exit 0
    fi

    getadempierestatus
    if [ $ADEMPIERESTATUS -eq 0 ] ; then
        echo "ADempiere is already running"
        return 1
    fi
    echo -n "Starting ADempiere ERP: "
    export LOGFILE=${ADEMPIERE_HOME}/jboss/server/adempiere/log/adempiere_`date +%Y%m%d%H%M%S`.log
    $SU ${ADEMPIERE_USER} -c "mkdir -p ${ADEMPIERE_HOME}/jboss/server/adempiere/log"
    $SU ${ADEMPIERE_USER} -c "cd ${ADEMPIERE_HOME}/utils;${ADEMPIERE_HOME}/utils/RUN_Server2.sh &> $LOGFILE &"
    RETVAL=$?
    if [ $RETVAL -eq 0 ] ; then
        # wait for server to be confirmed as started in logfile
        STATUSTEST=0
        ITERATIONS=0
        while [ $STATUSTEST -eq 0 ] ; do
            sleep $SLEEPSECONDS
            tail -n 9 $LOGFILE | grep -q 'INFO.*\[Server\].*Started in' && STATUSTEST=1
            echo -n "."
            ITERATIONS=`expr $ITERATIONS + 1`
            if [ $ITERATIONS -gt $MAXITERATIONS ]
                then
                break
            fi
        done
        if [ $STATUSTEST -eq 0 ]
        then
            log_warning_msg "Service hasn't started within the timeout allowed, please review file $LOGFILE to see the status of the service"
        else
            log_success_msg "Service started"
        fi
        echo
    else
        log_failure_msg "Service not started"
    echo
    fi
    return $RETVAL
}

stop () {

    if [ "${CONFIGURE_RUN}" != "true" ]		
    then
	echo "ADempiere is not configured.  You must run
'/etc/init.d/adempiere configure' as the root user to configure the server."
	exit 0
    fi

    getadempierestatus
    if [ $ADEMPIERESTATUS -ne 0 ] ; then
	echo "ADempiere is already stopped"
	return 1
    fi
    echo -n "Stopping ADempiere ERP: "
    # export LASTLOG=`ls -t ${ADEMPIERE_HOME}/jboss/server/adempiere/log/adempiere_??????????????.log | head -1`
    export LASTLOG=${ADEMPIERE_HOME}/jboss/server/adempiere/log/server.log
    $SU ${ADEMPIERE_USER} -c "cd ${ADEMPIERE_HOME}/utils;${ADEMPIERE_HOME}/utils/RUN_Server2Stop.sh &> /dev/null &"
    RETVAL=$?
    if [ $RETVAL -eq 0 ] ; then
	# wait for server to be confirmed as halted in logfile
	STATUSTEST=0
	ITERATIONS=0
	while [ $STATUSTEST -eq 0 ] ; do
	    sleep $SLEEPSECONDS
	    tail -n 9 $LASTLOG | grep -q "$STOPMESSAGE" && STATUSTEST=1
	    echo -n "."
	    ITERATIONS=`expr $ITERATIONS + 1`
	    if [ $ITERATIONS -gt $MAXITERATIONS ]
	    then
		break
	    fi
	done
	if [ $STATUSTEST -eq 0 ]
	then
	    log_warning_msg "Service hasn't stopped within the timeout allowed, please review file $LASTLOG to see the status of the service"
	    log_warning_msg "Trying direct kill with signal -15"
	    # Adempiere didn't finish - try direct kill with signal 15, then signal 9
	    kill -15 `ps ax | grep -v grep | grep ${ADEMPIERE_HOME} | sed -e 's/^ *//g' | cut -f 1 -d " "`
	    sleep 5
	    getadempierestatus
	    if [ $ADEMPIERESTATUS -ne 0 ] ; then
		log_success_msg "Service stopped with kill -15"
	    else
		echo "Trying direct kill with signal -9"
		kill -9 `ps ax | grep -v grep | grep ${ADEMPIERE_HOME} | sed -e 's/^ *//g' | cut -f 1 -d " "`
		sleep 5
		getadempierestatus
		if [ $ADEMPIERESTATUS -ne 0 ] ; then
		  log_success_msg "Service stopped with kill -9"
		else
		  log_warning_msg "Service hasn't stopped"
		fi
	    fi
	else
	    log_success_msg "Service stopped"
	fi
	echo
    else
	log_failure_msg "Service not stopped"
	echo
    fi
    return $RETVAL
}

restart () {
    stop
    sleep $SLEEPSECONDS
    start
}

condrestart () {
    getadempierestatus
    if [ $ADEMPIERESTATUS -eq 0 ] ; then
	restart
    fi
}

dostatus () {
    getadempierestatus
    if [ $ADEMPIERESTATUS -eq 0 ] ; then
	echo
	echo "ADempiere is running:"
	ps ax | grep -v grep | grep ${ADEMPIERE_HOME} | sed 's/^[[:space:]]*\([[:digit:]]*\).*:[[:digit:]][[:digit:]][[:space:]]\(.*\)/\1 \2/'
	echo
    else
	echo "ADempiere is stopped"
    fi
}

case "$1" in
    start)
	if test -f "$CONFIGURATION"
	then
            if test "${RUN_AT_STARTUP}" != "true" 
            then
		exit 0
	    fi
	else
	    echo "ADempiere is not configured.  You must run
'/etc/init.d/adempiere configure' as the root user to configure the server."
	    exit 0
	fi
	start
	;;
    configure)
        configure
        ;;
    stop)
	stop
	;;
    restart|reload|force-reload)
	restart
	;;
    condrestart)
	condrestart
	;;
    status)
	dostatus
	;;
    enable) 
	if test -f "$CONFIGURATION"
	then
            RUN_AT_STARTUP=true
	    write_sysconfig
	else
	    echo "ADempiere is not configured.  You must run
'/etc/init.d/adempiere configure' as the root user to configure the server."
	    exit 0
	fi
        ;;
  disable)
	if test -f "$CONFIGURATION"
	then
            RUN_AT_STARTUP=false
	    write_sysconfig
	else
	    echo "ADempiere is not configured.  You must run
'/etc/init.d/adempiere configure' as the root user to configure the server."
	    exit 0
	fi
	;;
    *)
        echo $"Usage: $0 {start|stop|restart|force-reload|configure|status|enable|disable}"
        exit 1

esac

exit 0
